# ld.so.preload
/home/user/MyGit/PPHook/build/use/libhook1.so
/home/user/MyGit/PPHook/build/use/libhook2.so
