# ld.so.preload
/home/user/MyGit/PerfectProloadHook/build/use/libhook1.so
/home/user/MyGit/PerfectProloadHook/build/use/libhook2.so
